# cloudflared.service
[Unit]
Description=cloudflared tunnel (conjure → conjure.vlad-p.com)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The tunnel token is fetched once with:
#   GET /accounts/{acct}/cfd_tunnel/{tunnel-id}/token
# and stored here (chmod 600, root-owned). It is NOT committed to git.
EnvironmentFile=/etc/cloudflared/conjure.env
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate run --token ${TUNNEL_TOKEN}
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
